# alt-shift-notify/alt-shift-notify.service.in
[Unit]
Description=Alt+Shift Toggle Notifier
After=local-fs.target
After=systemd-udevd.service
Requires=systemd-udevd.service
Requires=local-fs.target

[Service]
ExecStart="@INSTALL_DIR@/venv/bin/python" main.py
WorkingDirectory=@INSTALL_DIR@
Type=simple

[Install]
WantedBy=user-1000.slice

